# systemd/system/mnt-exthdd.automount
[Unit]
Description=Automount /mnt/exthdd

[Automount]
Where=/mnt/exthdd

[Install]
WantedBy=multi-user.target
